# templates/cdb-blocklist-fetch.timer.j2
[Unit]
Description=Daily blocklist fetch

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target

